# apps/supervisor/fhem.conf
[program:fhem]
directory = /opt/fhem
command = /opt/bin/fhem.start
user = root
autorestart = true
