# apps/supervisor/node-red.conf
[program:node-red]
directory = /opt/node-red
command = node red.js
user = root
