# apps/supervisor/node-red.conf
[program:node-red]
directory = /opt/node-red
command = node --max-old-space-size=128 red.js
user = root
environment = HOME="/root",USER="root"
autorestart = true
