# apps/supervisor/rfxtrx.conf
[program:rfxtrx remote serial port]
directory = /tmp
command = /opt/bin/socat.remotetrx.start
user = root
autorestart = true
