# apps/supervisor/conf.d.dist/emonhub.conf
[program:emonhub]
directory = /tmp
command = python /usr/share/emonhub/emonhub.py --config-file /etc/emonhub/emonhub.conf
user = root
autorestart = true
