# apps/supervisor/conf.d.dist/hdmi-cec-node-test.conf
[program:hdmi-cec-node-test]
directory = /tmp
command = /usr/local/bin/node /opt/hdmi-cec-node-test/bin/www
user = root
autorestart = true
