# apps/supervisor/conf.d.dist/mount-rwdata.conf
[program:mount-rwdata]
directory = /opt/bin
command = /opt/bin/mount-rwdata
user = root
autorestart = true
