# apps/supervisor/conf.d.dist/node-red-pi.conf
[program:node-red]
directory = /opt/node-red
command = node-red-pi --max-old-space-size=128
user = root
environment = HOME="/root",USER="root"
autorestart = true
