# apps/supervisor/conf.d.dist/openzwave-control-panel.conf
[program:openzwave-control-panel]
directory = /opt/ozwcp
command = /opt/ozwcp/ozwcp -p 8000
user = root
autorestart = true
