# apps/supervisor/conf.d.dist/ser2net.conf
[program:ser2net]
directory = /tmp
command = /usr/sbin/ser2net -c /etc/ser2net.conf -n
user = root
autorestart = true
