# apps/supervisor/conf.d.dist/socat.zwave.conf
[program:zwave remote serial port]
directory = /tmp
command = /opt/bin/socat.remotezwave.start
user = root
autorestart = true
