# apps/supervisor/conf.d.dist/weewx.conf
[program:weewx]
directory = /tmp
command = /usr/bin/python /usr/bin/weewxd /etc/weewx/weewx.conf
user = root
autorestart = true
